# cowrie_honeyfs/etc/systemd/system/app.target
[Unit]
Description=App Stack
Wants=app-web.service
